# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/env sh

exec /usr/sbin/php-fpm7 -c /etc/php/php-fpm.conf -F -R
