# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/env sh

exec php-fpm -c /etc/php/php-fpm.conf -F -R
